# .config/systemd/user/backup-config.timer
[Unit]
Description=Backup config files weekly

[Timer]
OnActiveSec=1min
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target

