# .config/systemd/user/backup-config.timer
[Unit]
Description=Backup Config everyday

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target

