# .config/systemd/user/backup-config.timer
[Unit]
Description=Backup config files daily

[Timer]
OnActiveSec=5min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
