# .config/systemd/user/backup-config.timer
[Unit]
Description=Backup config files weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

