# .config/systemd/user/backup-config.timer
[Unit]
Description=Backup config files 5 mins after boot

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
