# .config/systemd/user/backup-config.timer
[Unit]
Description=Backup config files weekly

[Timer]
# OnCalendar=weekly
OnBootSec=10s
Persistent=true

[Install]
WantedBy=timers.target

