# .config/systemd/user/data-backup.service
[Unit]
Description=Backup Data

[Service]
ExecStart=/home/aome510/scripts/snapbackup
