# autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=xbindkeys
Exec=xbindkeys
Comment[en_US]=Allow xbindkeys to run on startup
